# internal/doctor/testdata/fixtures/persistence/cron/daily-summary.cron
0 8 * * * openclaw run daily-client-email-summary
